# ansible/roles/trapped-mind/files/trapped-mind-update.service
[Unit]
Description=Check for trapped-mind updates and apply

[Service]
Type=oneshot
ExecStart=/usr/local/bin/trapped-mind-update.sh
