# ansible/roles/trapped-mind/files/trapped-mind-update.timer
[Unit]
Description=Periodically check for trapped-mind updates

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
